# xilinx-zc702-2025.2/project-spec/meta-user/recipes-modules/mydriver/files/myinit.service
[Unit]
Description=My Init Script
After=network.target

[Service]
ExecStart=/usr/bin/myinit.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target